# TimeShift_Automate/new_timeshift-8h.timer

[Unit]
Description=Run Timeshift every 6 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target


