# TimeShift_Automate/new_timeshift-boot.service
[Unit]
Description=Run Timeshift snapshot at boot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/New_timeshift-job.sh Boottime

