# TimeShift_Automate/new_timeshift-8h.service
[Unit]
Description=Run Timeshift snapshot every 8 hours

[Service]
Type=oneshot
ExecStart=/usr/local/bin/New_timeshift-job.sh Other

