# TimeShift_Automate/new_timeshift-8h.timer

[Unit]
Description=Run Timeshift every 12 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h
Persistent=true

[Install]
WantedBy=timers.target


