# TimeShift_Automate/timeshift-cleanup.service

[Unit]
Description=Timeshift Cleanup Script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleanup_timeshift.sh


